# container/etc/supervisord-noauth.conf
; A good starting point with documentation for this file can be generated
; from a supervisord installation by running echo_supervisord_conf

[supervisord]
logfile=/var/log/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=debug
pidfile=/var/run/supervisord.pid
nodaemon=true
minfds=1024
minprocs=200

[unix_http_server]
file=/var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:vouch]
command=/usr/local/bin/vouch --config /etc/vouch/vouch-noauth.conf
autostart=true
autorestart=true
startretries=99
startsecs=1
redirect_stderr=true
stdout_logfile=/var/log/vouch.log
user=root

[program:confd]
command=/usr/bin/confd -backend consul -node %(ENV_CONFIG_HOST_AND_PORT)s -scheme %(ENV_CONFIG_SCHEME)s -watch -prefix /customers/%(ENV_CUSTOMER_ID)s -auth-token %(ENV_CONSUL_HTTP_TOKEN)s
redirect_stderr=true
stdout_logfile=/var/log/confd/confd.log
autorestart=true
startretries=99
startsecs=1
